# log-server/platform/solaris/smf-manifest.xml.in
<?xml version='1.0'?>

<!--
# Copyright 2013-2014 Will Mason
# 
#    Licensed under the Apache License, Version 2.0 (the "License");
#    you may not use this file except in compliance with the License.
#    You may obtain a copy of the License at
# 
#      http://www.apache.org/licenses/LICENSE-2.0
# 
#    Unless required by applicable law or agreed to in writing, software
#    distributed under the License is distributed on an "AS IS" BASIS,
#    WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
#    See the License for the specific language governing permissions and
#    limitations under the License.
-->

<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='chucho'>
    <service name='application/chucho' type='service' version='1'>
        <create_default_instance enabled='true'/>
        <single_instance/>
        <dependency name='network'
                    grouping='require_all'
                    restart_on='none'
                    type='service'>
            <service_fmri value='svc:/milestone/network:default'/>
        </dependency>
        <exec_method type='method'
                     name='start'
                     exec='${CMAKE_INSTALL_PREFIX}/bin/chuchod --service'
                     timeout_seconds='60'/>
        <exec_method type='method'
                     name='stop'
                     exec=':kill'
                     timeout_seconds='60'/>
        <stability value='Evolving'/>
        <template>
            <common_name>
                <loctext xml:lang='C'>Chucho Log ${CHUCHO_VERSION}</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
