# log-server/platform/linux/systemd.service.in
#
# Copyright 2013-2016 Will Mason
# 
#    Licensed under the Apache License, Version 2.0 (the "License");
#    you may not use this file except in compliance with the License.
#    You may obtain a copy of the License at
# 
#      http://www.apache.org/licenses/LICENSE-2.0
# 
#    Unless required by applicable law or agreed to in writing, software
#    distributed under the License is distributed on an "AS IS" BASIS,
#    WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
#    See the License for the specific language governing permissions and
#    limitations under the License.
#

[Unit]
Description=Chucho logging framework log server
After=network.target
Requires=network.target

[Service]
# When running under systemd, don't use the --service argument. This
# would cause the pid file not to get created, which is how systemd
# tracks the actual pid.
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/chuchod
Type=forking
PIDFile=/var/run/chuchod.pid

[Install]
WantedBy=default.target
